# scripts/github-actions-runner.service
[Unit]
Description=GitHub Actions Runner (home-server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/actions-runner
ExecStart=/usr/bin/sg docker -c '%h/actions-runner/run.sh'
Restart=always
RestartSec=5
KillMode=process
TimeoutStopSec=30

[Install]
WantedBy=default.target
